# etc/udev/rules.d/99-update-refresh-rate.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/update-monitors.sh"
